# config/systemd/user/dwmstatus.timer
[Unit]
Description=Run dwmstatus every minute

[Timer]
OnActiveSec=0
OnCalendar=*:*:0/10
AccuracySec=1
