# config/systemd/user/dwmstatus.service
[Unit]
Description=Update dwm status text

[Service]
Type=oneshot
ExecStart=/usr/local/bin/dwmstatus
